# service/systemd/user/kairpodsd.service
[Unit]
Description=kAirPods D-Bus Service
After=graphical-session.target

[Service]
Type=dbus
BusName=org.kairpods
ExecStart=/usr/bin/kairpodsd
Restart=on-failure
RestartSec=5
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=yes
ReadWritePaths=/var/run/dbus

[Install]
WantedBy=default.target
